# overlay-voguev210/sys-apps/voguev210-rtc-war/files/date-start.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Sets the time from a value stored on disk (if system is older)."
author        "chromium-os-dev@chromium.org"

start on stopping startup

task

script

current_time=`date "+%s"`
SAVED_TIME_FILE=/mnt/stateful_partition/etc/saved_time
if [ -f $SAVED_TIME_FILE ] ; then
  saved_time=`cat $SAVED_TIME_FILE`
else
  saved_time=`date -d "2010/1/1 UTC" "+%s"`
fi

if [ $current_time -lt $saved_time ] ; then
  date --set "@$saved_time"
fi

end script

