# overlay-voguev210/sys-apps/voguev210-mac/files/voguev210-mac.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Vogue V210 eth0 MAC setup"
author        "chromium-os-dev@chromium.org"

start on starting connman

task

script

mac=$(uboot-env.py -f /dev/mtdblock0 -o 0x7c000 -s 0x4000 --get ethaddr)
if [ -z "$mac" ] ; then
  echo "Invalid MAC, run setup_mac"
  exit 1
fi
if [ "$mac" = "00:40:5c:26:0a:5b" ] ; then
  echo "WARNING: you are using the default MAC, you need to run setup_mac"
fi

ifconfig eth0 hw ether $mac

end script

